# pimame_menu.conf
[program:pimame_menu]
directory=/home/pi/pimame/pimame-menu/
command=python launchmenu.py
autostart=true
stderr_logfile=/var/log/long.err.log
stdout_logfile=/var/log/long.out.log
